# ansible/roles/backup/templates/180dc-backup.service.j2
[Unit]
Description=180DC database, runtime volume, and log backup
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
Environment=RETENTION_COUNT={{ backup_retention_count }}
ExecStart={{ backup_root }}/backup.sh
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
